# deploy/urpstudio-backend.service
[Unit]
Description=urpstudio backend API service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/urpstudio
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/usr/bin/node /opt/urpstudio/server/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
